# giga-nepal-backend/deploy/systemd/neogiga-queue-marketing.service
[Unit]
Description=NeoGiga campaign preparation and marketing queue worker
After=network-online.target postgresql.service

[Service]
Type=simple
User=neogiga
Group=neogiga
WorkingDirectory=/home/neogiga/laravel/current
Nice=10
CPUWeight=20
UMask=0002
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=360
ExecStart=/usr/bin/php /home/neogiga/laravel/current/artisan queue:work database --queue=campaign-preparation,marketing --sleep=3 --tries=5 --timeout=120 --max-time=3600 --memory=256
StandardOutput=append:/home/neogiga/laravel/current/storage/logs/marketing-queue-worker.log
StandardError=append:/home/neogiga/laravel/current/storage/logs/marketing-queue-worker.log

[Install]
WantedBy=multi-user.target
