# giga-nepal-backend/deploy/systemd/neogiga-queue-default.service
[Unit]
Description=NeoGiga default queue worker (scheduled jobs, core operations)
After=network-online.target postgresql.service

[Service]
Type=simple
User=neogiga
Group=neogiga
WorkingDirectory=/home/neogiga/laravel/current
Nice=5
CPUWeight=50
UMask=0002
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=360
ExecStart=/usr/bin/php /home/neogiga/laravel/current/artisan queue:work database --queue=default --sleep=3 --tries=3 --timeout=120 --max-time=3600 --memory=256
StandardOutput=append:/home/neogiga/laravel/current/storage/logs/default-queue-worker.log
StandardError=append:/home/neogiga/laravel/current/storage/logs/default-queue-worker.log

[Install]
WantedBy=multi-user.target
