# giga-nepal-backend/deploy/systemd/neogiga-queue-imports.service
[Unit]
Description=NeoGiga customer and catalog import queue worker
After=network-online.target postgresql.service

[Service]
Type=simple
User=neogiga
Group=neogiga
WorkingDirectory=/home/neogiga/laravel/current
Nice=10
CPUWeight=20
UMask=0002
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=1860
ExecStart=/usr/bin/php /home/neogiga/laravel/current/artisan queue:work database_imports --queue=imports,catalog-imports --sleep=3 --tries=3 --timeout=1800 --max-time=3600 --memory=512
StandardOutput=append:/home/neogiga/laravel/current/storage/logs/imports-queue-worker.log
StandardError=append:/home/neogiga/laravel/current/storage/logs/imports-queue-worker.log

[Install]
WantedBy=multi-user.target
